# system/urms.cron
* * * * * /home/sr99_ramp_metering/system/check_and_restart_urms.sh
